# scripts/systemd/clawops-release-watch.timer
[Unit]
Description=Run clawops release watcher hourly

[Timer]
OnBootSec=10min
OnUnitActiveSec=1h
RandomizedDelaySec=5min
Persistent=true

[Install]
WantedBy=timers.target
