# scripts/systemd/clawops-daily-digest.service
[Unit]
Description=Clawops Daily Release Digest
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/clawops
SyslogIdentifier=clawops-daily-digest
ExecStart=/usr/bin/sudo -n /usr/bin/docker exec clawops-gateway sh -lc "set -e; set -a; [ -f /opt/clawops/.env ] && . /opt/clawops/.env; set +a; /opt/clawops/scripts/daily-release-digest.sh"
