# scripts/systemd/clawops-weekly-hygiene.timer
[Unit]
Description=Run clawops weekly hygiene

[Timer]
OnCalendar=Sun *-*-* 13:00:00
RandomizedDelaySec=20min
Persistent=true

[Install]
WantedBy=timers.target
